# config/supervisord/conf.d/credentials.conf
[program:credentials]
command=/generate-credentials.sh
